# roles/nextcloud_config/files/nextcloud-optimize.timer
[Unit]
Description=Run Nextcloud Optimize Script Daily @ 01:15

[Timer]
OnCalendar=*-*-* 1:15:00
Persistent=true 

[Install]
WantedBy=timers.target